# puppet/modules/mopidy/files/mopidy-monitor.timer
[Unit]
Description=Mopidy Check Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m

[Install]
WantedBy=basic.target
