# puppet/modules/mopidy/templates/mopidy-monitor.service.erb
[Unit]
Description=Mopidy Check service
After=mopidy.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mopidy-cli resume-radio
User=mopidy
